# dnmp/services/apache/supervisord.conf
[supervisord]
nodaemon=true

[program:apache]
command=/usr/local/apache/bin/httpd
